# wg-agent-watchdog.service
# One-shot service run by wg-agent-watchdog.timer to check health and restart wg-agent if needed.

[Unit]
Description=wg-agent health check (watchdog)
After=wg-agent.service

[Service]
Type=oneshot
ExecStart=INSTALL_DIR/check-wg-agent-health.sh
